# doc/example/eneverre.service
# Example systemd unit for the Eneverre NVR API.
#
# Install:
#   sudo install -m0755 eneverre /usr/local/bin/eneverre
#   sudo install -d /etc/eneverre/cameras.d
#   sudo install -m0644 doc/example/eneverre.ini /etc/eneverre/eneverre.ini
#   sudo cp doc/example/cameras.d/*.ini /etc/eneverre/cameras.d/
#   sudo install -m0644 doc/example/eneverre.service /etc/systemd/system/eneverre.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now eneverre
#
# Then check it:  systemctl status eneverre   /   journalctl -u eneverre -f
#
# Config is read from /etc/eneverre/eneverre.ini and /etc/eneverre/cameras.d/.
# Persistent state is a single SQLite DB (it also holds the rotating stream
# credentials) in /var/lib/eneverre, created and owned automatically via
# StateDirectory.
#
# EMBEDDED MEDIA ENGINE ([media] in eneverre.ini): recordings + the segment
# index default to /var/lib/eneverre/recordings, which is under StateDirectory
# and therefore writable despite ProtectSystem=strict — no change needed. If
# you point [media] record_dir at a path OUTSIDE /var/lib/eneverre, add it to
# ReadWritePaths (e.g. `ReadWritePaths=/srv/recordings`) or writes will fail.
# The engine also opens an RTSP relay (default :8554); make sure that port is
# reachable by your RTSP clients (it does not go through the reverse proxy).
#
# NOTE on file modes: DynamicUser=yes runs the service as a transient UID/GID,
# so config files must stay readable by it (0644, as above). If eneverre.ini
# holds secrets ([events] webhook_secret, [media] record paths) on a multi-user
# host, switch to a dedicated account instead: drop `DynamicUser=yes`, create
# `eneverre` (e.g. `useradd --system --no-create-home --shell /usr/sbin/nologin
# eneverre`), set `User=eneverre` / `Group=eneverre`, and
# `chown root:eneverre /etc/eneverre/eneverre.ini && chmod 0640` it.

[Unit]
Description=Eneverre NVR API
Documentation=https://github.com/matiasdelellis/eneverre-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/eneverre
Restart=on-failure
RestartSec=2

# Run as a transient, isolated system user (no useradd needed). systemd creates
# /var/lib/eneverre, owns it as that user, and keeps it across restarts.
DynamicUser=yes
StateDirectory=eneverre
WorkingDirectory=/var/lib/eneverre
Environment=ENEVERRE_DB_PATH=/var/lib/eneverre/eneverre.db

# Admin user. On first start (empty users table) the service creates an admin
# with a RANDOM password and logs it once — read it with
# `journalctl -u eneverre | grep 'generated password'`, then log in and change
# it. No password is stored in this unit. For unattended provisioning you may
# set a known one via a drop-in (`systemctl edit eneverre`):
#   [Service]
#   Environment=ENEVERRE_ADMIN_USER=admin
#   Environment=ENEVERRE_ADMIN_PASS=your-secret
# Either way it is used only while the users table is empty.

# Log level (debug | info | warn | error). The access log is one line/request.
Environment=ENEVERRE_LOG_LEVEL=info

# --- Hardening ------------------------------------------------------------
# The default listen port is 8080. To bind a privileged port (< 1024), add:
#   AmbientCapabilities=CAP_NET_BIND_SERVICE
#   CapabilityBoundingSet=CAP_NET_BIND_SERVICE
# otherwise the empty CapabilityBoundingSet below drops all capabilities.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
